# docs/handbook/server/supervisord/conf.d/test.conf
[program:test]
command=/srv/scripts/61/supervisor/test.py youja 
autostart=false
autorestart=true
user=jesse
redirect_stderr=true
stdout_logfile=/dev/shm/test.out
stderr_logfile=/dev/shm/test.err
