# docs/handbook/server/os/ubuntu/profile.d/ffmpeg-env.sh
export FFMPEG_HOME=/usr/local/ffmpeg
export PATH=${FFMPEG_HOME}/bin:$PATH

