# docs/handbook/server/os/ubuntu/profile.d/go-env.sh
export GO_HOME=/usr/local/go
export PATH=${GO_HOME}/bin:$PATH
export GOROOT=${GO_HOME}
export GOPATH=${GO_HOME}
export GOBIN=${GO_HOME}/bin
