# docs/handbook/server/os/ubuntu/profile.d/iperf-env.sh
export IPERF_HOME=/usr/local/iperf
export PATH=${IPERF_HOME}/bin:$PATH

