# docs/handbook/server/os/ubuntu/profile.d/maven-env.sh
export M2_HOME=/usr/local/maven
export CLASSPATH=${M2_HOME}:$CLASSPATH
export PATH=${M2_HOME}/bin:$PATH

