# docs/handbook/server/os/ubuntu/profile.d/python-env.sh
export PYTHON_HOME=/usr/local/python
export PATH=${PYTHON_HOME}/bin:$PATH
