# docs/handbook/server/os/mac/profile.d/1-env.sh
export PATH=$PATH:/usr/local/sbin:~/bin/ext/mac:~/bin/ext/gat
