# docs/handbook/server/os/mac/profile.d/alias.sh

#base
alias ll='ls -l'

alias cdwork='cd ~/workspace'
alias cdbin='cd ~/bin'
alias cddoc='cd ~/Documents'
alias cddown='cd ~/Downloads'

alias cdlocal='cd /usr/local'
alias cdredis='cd /usr/local/redis'
alias cdmongo='cd /usr/local/mongodb'
alias cdnginx='cd /usr/local/Cellar/nginx'
alias cdnginxconf='cd /usr/local/etc/nginx'


#youja.cn
alias cdywork='cd ~/workspace/youja.cn'
alias cdydns='cd ~/workspace/youja.cn/yunwei/server-config/dnsmasq'
alias cdysvn='cd ~/workspace/youja.cn/yunwei'

#518yin.com
alias cdhwork='cd ~/workspace/518yin.com'
alias cdhsvn='cd ~/workspace/518yin.com/yunwei'

#8ops.cc
alias cdjwork='cd ~/workspace/8ops.top'
alias cdjupx='cd ~/workspace/8ops.top/_upx'
alias cdjenv='cd ~/workspace/8ops.top/_env'

#guanaitong.com
alias cdgwork='cd ~/workspace/guanaitong.com'
alias cdgsvn='cd ~/workspace/guanaitong.com/sysom'
alias cdgo='cd ~/workspace/go/src'

#autobestdevops.com
alias cdawork='cd ~/workspace/autobest.com'
alias cdahandbook='cd ~/workspace/autobest.com/handbook'



