# docs/handbook/server/os/mac/profile.d/clojure-env.sh
export CLOJURE_HOME=/usr/local/clojure
export PATH=${CLOJURE_HOME}:${PATH}
