# docs/handbook/server/os/mac/profile.d/go-env.sh
export GO_HOME=/usr/local/go
export PATH=${GO_HOME}/bin:$PATH
export GOROOT=${GO_HOME}
#export GOBIN=${GO_HOME}/bin
export GOPATH=~/workspace/go
export PATH=${GOPATH}/bin:$PATH
#export GO111MODULE=on
