# docs/handbook/server/os/mac/profile.d/libxml2-env.sh
export LIBXML2_HOME=/usr/local/opt/libxml2 
export PATH=${LIBXML2_HOME}/bin:$PATH
