# docs/handbook/server/os/mac/profile.d/mongo-env.sh
export MONGODB_HOME=/usr/local/mongodb
export PATH=${MONGODB_HOME}/bin:$PATH
