# docs/handbook/server/os/mac/profile.d/phantomjs-env.sh
export PHANTOMJS_HOME=/usr/local/phantomjs
export PATH=${PHANTOMJS_HOME}/bin:$PATH
