# docs/handbook/server/os/mac/profile.d/ruby-env.sh
export RUBY_HOME=/usr/local/ruby
export PATH=${RUBY_HOME}/bin:$PATH
