# docs/handbook/server/os/mac/profile.d/sqlite3-env.sh
export SQLITE3_HOME=/usr/local/sqlite3
export PATH=${SQLITE3_HOME}/bin:$PATH
