# docs/handbook/server/os/mac/profile.d/tomcat-env.sh
export TOMCAT_HOME=/usr/local/tomcat
export CATALINA_HOME=${TOMCAT_HOME}
export PATH=${TOMCAT_HOME}/bin:$PATH
